# os-reticulum/pkg-plist
/usr/local/etc/rc.d/rnsd
/usr/local/etc/rc.d/lxmd
/usr/local/etc/rc.syshook.d/start/50-reticulum
/usr/local/etc/newsyslog.conf.d/reticulum.conf
/usr/local/opnsense/mvc/app/models/OPNsense/Reticulum/ACL/ACL.xml
/usr/local/opnsense/mvc/app/models/OPNsense/Reticulum/Menu/Menu.xml
/usr/local/opnsense/mvc/app/models/OPNsense/Reticulum/Reticulum.xml
/usr/local/opnsense/mvc/app/models/OPNsense/Reticulum/Reticulum.php
/usr/local/opnsense/mvc/app/controllers/OPNsense/Reticulum/IndexController.php
/usr/local/opnsense/mvc/app/controllers/OPNsense/Reticulum/Api/RnsdController.php
/usr/local/opnsense/mvc/app/controllers/OPNsense/Reticulum/Api/LxmdController.php
/usr/local/opnsense/mvc/app/controllers/OPNsense/Reticulum/Api/ServiceController.php
/usr/local/opnsense/mvc/app/views/OPNsense/Reticulum/general.volt
/usr/local/opnsense/mvc/app/views/OPNsense/Reticulum/interfaces.volt
/usr/local/opnsense/mvc/app/views/OPNsense/Reticulum/lxmf.volt
/usr/local/opnsense/mvc/app/views/OPNsense/Reticulum/logs.volt
/usr/local/opnsense/scripts/OPNsense/Reticulum/reconfigure.sh
/usr/local/opnsense/scripts/OPNsense/Reticulum/rnsd_status.sh
/usr/local/opnsense/scripts/OPNsense/Reticulum/lxmd_status.sh
/usr/local/opnsense/scripts/OPNsense/Reticulum/rnstatus.sh
/usr/local/opnsense/scripts/OPNsense/Reticulum/info.sh
/usr/local/opnsense/service/conf/actions.d/actions_reticulum.conf
/usr/local/opnsense/service/templates/OPNsense/Reticulum/+TARGETS
/usr/local/opnsense/service/templates/OPNsense/Reticulum/reticulum_config.j2
/usr/local/opnsense/service/templates/OPNsense/Reticulum/lxmf_config.j2
/usr/local/opnsense/service/templates/OPNsense/Reticulum/rc.conf.d_rnsd.j2
/usr/local/opnsense/service/templates/OPNsense/Reticulum/rc.conf.d_lxmd.j2
/usr/local/opnsense/service/templates/OPNsense/Reticulum/lxmf_allowed.j2
/usr/local/opnsense/service/templates/OPNsense/Reticulum/lxmf_ignored.j2
/usr/local/opnsense/www/js/widgets/Reticulum.js
/usr/local/opnsense/www/js/widgets/Metadata/Reticulum.xml
