# os-reticulum/src/etc/rc.d/rnsd
#!/bin/sh

# PROVIDE: rnsd
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="rnsd"
rcvar="rnsd_enable"

load_rc_config $name

: ${rnsd_enable:="NO"}
: ${rnsd_user:="reticulum"}
: ${rnsd_config:="/usr/local/etc/reticulum"}

pidfile="/var/run/${name}.pid"
command="/usr/local/reticulum-venv/bin/rnsd"
command_args="--config ${rnsd_config}"
command_interpreter="/usr/local/reticulum-venv/bin/python3"
command_user="${rnsd_user}"
rnsd_flags=""

start_precmd="${name}_prestart"
stop_postcmd="${name}_poststop"

rnsd_prestart()
{
    # Ensure config directory exists with restrictive permissions (identity keys inside)
    if [ ! -d "${rnsd_config}" ]; then
        mkdir -p "${rnsd_config}"
    fi
    chown ${rnsd_user}:${rnsd_user} "${rnsd_config}"
    chmod 700 "${rnsd_config}"

    # Ensure log directory exists
    if [ ! -d "/var/log/reticulum" ]; then
        mkdir -p /var/log/reticulum
        chown ${rnsd_user}:${rnsd_user} /var/log/reticulum
    fi

    # Ensure identity key files created by rnsd are not world-readable
    umask 077
}

rnsd_poststop()
{
    rm -f ${pidfile}
}

run_rc_command "$1"
